# ops/curiator.net/rathole-sietch-healthcheck.timer
[Unit]
Description=Periodically check curIAtor public gallery tunnels

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
RandomizedDelaySec=15s
Persistent=true
Unit=rathole-sietch-healthcheck.service

[Install]
WantedBy=timers.target
